# machine2/etc/runit/sv/unison/run
#!/bin/sh
HOME=/home/user1 UNISON=/home/user1/.local/unison exec chpst -u user1 unison -socket 10000 2>&1
